# azure-api/supervisord.conf
[supervisord]
nodaemon=true
user=root

# [program:gunicorn]
# command=gunicorn --config gunicorn_config.py OmxServer:app

[program:nginx]
command=nginx -g "daemon off;"
